# Installation/arangodb.list
################################################################################
## epm configuration file
## Documentation: http://www.epmhome.org
################################################################################

$version=${arangodb_version}
$release=${arangodb_release}

################################################################################
## ArangoDB
################################################################################

%product ArangoDB
%copyright 2012 by triAGENS GmbH
%vendor triAGENS GmbH
%license COPYING
%readme README
%description ArangoDB Software
%version ${version}
%release ${release}

################################################################################
## File List
################################################################################

d 0755 ${susr} ${sgrp} /etc/arangodb -
d 0755 ${susr} ${sgrp} /etc/arangodb/ -
c 0644 ${susr} ${sgrp} /etc/arangodb/arangod.conf ${project_dir}/Installation/arangod.conf

f 0755 ${susr} ${sgrp} ${sbindir}/arangod-${arangodb_version} ${project_dir}/bin/arangod
l 000 ${rusr} ${rgrp} ${sbindir}/arangod ${sbindir}/arangod-${arangodb_version}

f 0755 ${susr} ${sgrp} ${bindir}/arangoimp ${project_dir}/bin/arangoimp
f 0755 ${susr} ${sgrp} ${bindir}/arangosh ${project_dir}/bin/arangosh
f 0755 ${susr} ${sgrp} ${bindir}/arango-password ${project_dir}/bin/arango-password
f 0755 ${susr} ${sgrp} ${bindir}/arango-dfdb ${project_dir}/bin/arango-dfdb

# database directory
d 0755 ${rusr} ${rgrp} ${data_dir} -
d 0755 ${susr} ${sgrp} ${data_dir}/arangodb -

# static files
d 0755 ${susr} ${sgrp} ${static_dir} -
d 0755 ${susr} ${sgrp} ${static_dir}/arangodb -
d 0755 ${susr} ${sgrp} ${static_dir}/arangodb/js -
d 0755 ${susr} ${sgrp} ${static_dir}/arangodb/html -

# PID file directory
d 0755 ${susr} ${sgrp} /var/run/arangodb -

# logfile directory
d 0755 ${susr} ${sgrp} /var/log/arangodb -

%include ${project_dir}/Installation/arangodb.sublist

################################################################################
## Pre/Post Install
################################################################################

## -----------------------------------------------------------------------------
## MAC OS X
## -----------------------------------------------------------------------------

%ifdef macosx
d 0755 ${rusr} ${rgrp} /Library/LaunchDaemons/ -
f 0755 ${rusr} ${rgrp} /Library/LaunchDaemons/org.arangodb.plist.template ${project_dir}/Installation/org.arangodb.plist
%postinstall <<EOF
sed -e 's%@BINARY@%${sbindir}/arangod-${version}%g' -e 's%@CONFIGDIR@%/etc/arangodb%g' /Library/LaunchDaemons/org.arangodb.plist.template > /Library/LaunchDaemons/org.arangodb.plist
launchctl unload /Library/LaunchDaemons/org.arangodb.plist
launchctl load /Library/LaunchDaemons/org.arangodb.plist
# we do not need this: launchctl start de.triagens.arango
EOF

%preremove <<EOF
launchctl stop de.triagens.arango
launchctl unload /Library/LaunchDaemons/org.arangodb.plist
EOF

## -----------------------------------------------------------------------------
## LINUX
## -----------------------------------------------------------------------------

%else
%preinstall <<EOF
getent group arango >/dev/null || groupadd -r arango 
getent passwd arango >/dev/null || useradd -r -g arango -d ${static_dir}/arangodb -s /bin/false -c "ArangoDB Application User" arango 
EOF
%endif

## -----------------------------------------------------------------------------
## CHKCONF
## -----------------------------------------------------------------------------

%ifdef chkconf
f 0755 ${rusr} ${rgrp} ${initdir}/arangod ${project_dir}/Installation/${START_SCRIPT}

%postinstall <<EOF
chkconfig --level ${runlevels} arangod on
exit 0
EOF

%preremove <<EOF
${initdir}/arangod stop 2&>1 > /dev/null || true
chkconfig --del arangod
exit 0
EOF

## -----------------------------------------------------------------------------
## INSSERV
## -----------------------------------------------------------------------------

%elseifdef insserv
f 0755 ${rusr} ${rgrp} ${initdir}/arangod ${project_dir}/Installation/${START_SCRIPT}

%postinstall <<EOF
test -x /sbin/insserv && /sbin/insserv /etc/init.d/arangod 2&>1 > /dev/null
exit 0
EOF

%preremove <<EOF
${initdir}/arangod stop 2&>1 > /dev/null || true
test -x /sbin/insserv && /sbin/insserv -r /etc/init.d/arangod 2&>1 > /dev/null
exit 0
EOF

## -----------------------------------------------------------------------------
## EPM
## -----------------------------------------------------------------------------

%else

%system linux
i 0755 ${rusr} ${rgrp} arangod ${project_dir}/Installation/${START_SCRIPT} runlevel(${runlevels})

%endif

## -----------------------------------------------------------------------------
## RC skript
## -----------------------------------------------------------------------------

%system linux
l 000 ${rusr} ${rgrp} ${sbindir}/rcarangod /etc/init.d/arangod 
